# backend/supervisor.conf
[program:codemafia-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/codemafia/backend/artisan queue:work database --sleep=3 --tries=3 --max-jobs=500 --timeout=90
directory=/var/www/codemafia/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/codemafia-queue.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:codemafia-reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/codemafia/backend/artisan reverb:start
directory=/var/www/codemafia/backend
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/codemafia-reverb.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:codemafia-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/codemafia/backend/artisan schedule:work
directory=/var/www/codemafia/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/codemafia-scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
